# terraform/station/watershop/duckpond/pond@water-prod.timer
[Unit]
Description=DuckPond water-prod (every 10min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
