# terraform/station/watershop/duckpond/pond@water-staging.timer
[Unit]
Description=DuckPond water-staging (every 10min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
