# terraform/station/watershop/duckpond/pond@noyo-prod.timer
[Unit]
Description=DuckPond noyo-prod (every 30min)

[Timer]
OnBootSec=6min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
