# config/systemd/pond-selfmon@.service
[Unit]
Description=DuckPond self-monitoring %i (native)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
MemoryAccounting=yes
ExecStart=/home/jmacd/duckpond/config/scripts/run-selfmon.sh %i
Environment=HOME=/home/jmacd
